# code-server/rootfs/etc/s6-overlay/s6-rc.d/custom-scripts/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Code Server
# Runs user defined init scripts and installs provided services
# ==============================================================================

# Disable scripts if log level set to debug
 if bashio::debug ; then
        bashio::log.debug "Disabling custom scripts"
        bashio::exit.ok
fi

# Check for needed directories and copy templates if necessary
if ! bashio::fs.directory_exists "/config/custom-services" ; then
    cp -rf /etc/s6-overlay/s6-rc.d/user2 /config/custom-services
fi

# Timestamp for dynamic s6 db creation
stamp=$(s6-clock)

# The s6 scripts need to be run in background as the s6 service is locked during startup phase
function post {

# Compile s6-rc database
s6-rc-compile -b \
    "/run/s6/db-${stamp}" \
    \
    /etc/s6-overlay/s6-rc.d \
    /package/admin/s6-overlay/etc/s6-rc/sources \
    /config/custom-services

# Change live db to newly created database
s6-rc-update -b "/run/s6/db-${stamp}"

# Start services defined in custom-services bundle
s6-rc -b start custom-services

}

post &
