# code-server/rootfs/etc/s6-overlay/s6-rc.d/cont-init/up
/bin/bash -c /etc/s6-overlay/scripts/spin-up.sh