# code-server/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Code Server
# Runs the code server
# ==============================================================================
bashio::log.info 'Starting cron service...'
exec 2>&1
exec cron -f -L 15
